# cmd/cliprd/build/cliprd.service
[Unit]
Description=Cliprd Daemon

[Service]
ExecStartPre=/bin/sleep 3
ExecStart=/absolute/path/to/your/cliprd
Restart=on-success

[Install]
WantedBy=default.target
