# cmd/cliprd/build/cliprd.service
[Unit]
Description=Cliprd Daemon

[Service]
ExecStart=/absolute/path/to/your/cliprd
Restart=on-success

[Install]
WantedBy=multi-user.target
